# config/supervisor.conf
[program:chat]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/home/file/chat
autostart=true
autorestart=true
